# cdr_organize/cdr_organize.cron
# CDR Organizer - Process incoming files every 5 minutes
# Place this file in /etc/cron.d/cdr_organize

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * cdr_admin python3 /home/cdr_admin/PondCDRSuite/cdr_organize/cdr_organize.py /home/cdr_admin/CDRs/inbound/company_cdr /srv/cdr-office
*/5 * * * * cdr_admin python3 /home/cdr_admin/PondCDRSuite/cdr_organize/cdr_organize.py /home/cdr_admin/CDRs/inbound/company_lu /srv/cdr-office
