# cdr_organize/cdr_organize.cron
# CDR Organizer - Process incoming files every 5 minutes
# Place this file in /etc/cron.d/cdr_organize

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * {USER} python3 {PROJECT_DIR}/cdr_organize/cdr_organize.py {SOURCE_DIR_CDR} {DEST_DIR}
*/5 * * * * {USER} python3 {PROJECT_DIR}/cdr_organize/cdr_organize.py {SOURCE_DIR_LU} {DEST_DIR}
